# deploy/systemd/tower.service
[Unit]
Description=Tower Control Script
Wants=network-online.target
After=network-online.target pigpiod.service
Requires=pigpiod.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/condorcam-tower

# Prefer module execution now that pi_runtime is a package
ExecStart=/home/pi/condorcam-tower/.venv/bin/python -m pi_runtime.main

Restart=on-failure
RestartSec=2

# Better logging behavior in journald
Environment=PYTHONUNBUFFERED=1
SyslogIdentifier=tower

[Install]
WantedBy=multi-user.target
